# ruby-passenger-sigsci/3.3-jammy/supervisord/supervisord.conf
[supervisord]
nodaemon = true
silent = true
logfile = /dev/stderr
logfile_maxbytes = 0
pidfile = /var/run/supervisord.pid

# Derived images can drop extra [program:*] files here to run alongside the agent and nginx.
[include]
files = /usr/src/supervisord/conf.d/*.conf

[program:sigsci-agent]
command = /usr/local/bin/sigsci-agent-wrapper
autorestart = true
startsecs = 0
priority = 10
stdout_logfile = /dev/stderr
stderr_logfile = /dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:nginx]
command = sudo -E /usr/sbin/nginx -c /usr/src/nginx/nginx.conf
autorestart = true
stopsignal = QUIT
stopasgroup = true
killasgroup = true
stopwaitsecs = %(ENV_NGINX_STOPWAITSECS)s
priority = 20
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
