# node-pm2/24/supervisord/supervisord.conf
[supervisord]
nodaemon = true
logfile = /var/log/app/supervisord.log
logfile_maxbytes = 5MB
logfile_backups = 5
pidfile = /var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -c /usr/src/nginx/nginx.conf
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:pm2]
command = /usr/local/bin/pm2 start /usr/src/app/.pm2.config.json -s --no-daemon --only app
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
